# crontab
* * * * * /root/gitmind/shell.sh &>/tmp/error
